# server_config/meal-planner.service
[Unit]
Description=Meal Planner (Docker Compose)
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/home/admin/meal-planner
ExecStartPre=/usr/bin/docker compose -f /home/admin/meal-planner/docker-compose.yml pull
ExecStart=/usr/bin/docker compose -f /home/admin/meal-planner/docker-compose.yml up -d
ExecStop=/usr/bin/docker compose -f /home/admin/meal-planner/docker-compose.yml down
TimeoutStartSec=0
User=admin
Group=admin

[Install]
WantedBy=multi-user.target
